# letsencrypt/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command
30 3 * * *  root /opt/docker/run.sh > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
